# deploy/gpu/scripts/comfyui.service
[Unit]
Description=ChronoCanvas ComfyUI Server
After=network-online.target
Wants=network-online.target

[Service]
User=comfyui
WorkingDirectory=/opt/ComfyUI
ExecStart=/opt/ComfyUI/.venv/bin/python main.py --listen 0.0.0.0 --port 8188
Restart=always
RestartSec=5
Environment=CUDA_VISIBLE_DEVICES=0
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
